# infra/oci/services/pricewatch-embed.service
[Unit]
Description=Pricewatch internal query-embedding service (natural-language search)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/pricewatch/services
ExecStart=/opt/pricewatch/.venv/bin/uvicorn embed_server:app --host 127.0.0.1 --port 8301 --workers 1
Restart=always
RestartSec=5
# Same 1GB-RAM VM as pricewatch-web/pricewatch-bot - separate cgroup from
# pricewatch-web on purpose, so an onnxruntime memory spike here can't take
# down the SSR/image-proxy service. Model loads lazily on first request.
MemoryMax=350M

[Install]
WantedBy=multi-user.target
